# deploy/supervisor/devhelp.conf
; DevHelp (ArcBook) — Daphne ASGI server
;
; Install:
;   sudo cp deploy/supervisor/devhelp.conf /etc/supervisor/conf.d/devhelp.conf
;   sudo mkdir -p /var/log/supervisor
;   sudo supervisorctl reread
;   sudo supervisorctl update
;
; Manage:
;   sudo supervisorctl status devhelp
;   sudo supervisorctl start devhelp
;   sudo supervisorctl stop devhelp
;   sudo supervisorctl restart devhelp

[program:devhelp]
command=/home/python/Desktop/dev/DevHelp/.venv/bin/daphne config.asgi:application -b 0.0.0.0 -p 8000
directory=/home/python/Desktop/dev/DevHelp
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/devhelp_out.log
stderr_logfile=/var/log/supervisor/devhelp_error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
